# Cron/crontab
* * * * * ls >> /var/log/cron.log
# 마지막 줄은 지우면 안됩니다.
